# root/etc/uci-defaults/20_10-config
#!/bin/sh
# vim: set noexpandtab tabstop=4 shiftwidth=4 softtabstop=4 :

uci -q set 'network.lan.ipaddr=192.168.100.1'
uci -q set 'system.@system[0].hostname=OverTheBox'
